# cookbooks/xacct/templates/default/upstart-service.conf.erb
description <%=@service_name%>
author minhluan
setuid root
setgid root

env JAVA_HOME=<%=@java_home%>
start on <%=@depend_service%>
stop on runlevel [!2345]
expect fork
respawn

post-start script
    echo "begin start"
    output=`ps ax | grep '<%=@search_string%>'`
    set -- $output
    pid=$1
    jsvc=$5
    pattern="jsvc.exec"
	searchString="<%=@search_string%>"
	servicePathExec="<%=@service_home%>/bin/<%=@service_name%>.sh"
	servicePathPid="<%=@service_home%>/bin/<%=@service_name%>.pid"
	echo $jsvc
    if [ -f $servicePathPid ] && [$jsvc = $pattern ]; then
     	echo "kill and remove $output $jsvc $pid"
       	rm $servicePathPid
      	kill -9 $(ps ax|grep '$searchString')
        exec "$servicePathExec" start
    elif [ $jsvc = $pattern ] && [ ! -f $servicePathPid ]; then
	set +e
            echo "start kill process $jsvc $pid"
	kill -9 `ps ax|grep $searchString|awk '{print $1}'`
	exec "$servicePathExec" start
    elif [ $jsvc != $pattern ] && [ -f $servicePathPid ]; then
	echo "remove pid file $output $jsvc $pid"
            rm $servicePathPid
            exec "$servicePathExec" start
    else
	echo "default start"
            exec "$servicePathExec" start
    fi
    echo "end start"
end script

post-stop script
	echo "begin stop"
	exec "<%=@service_home%>/bin/<%=@service_name%>.sh" stop
end script
